# integration_tests/inetd/container/etc/xinetd.d/discard
# default: off
# description: An RFC 863 discard server.
# This is the tcp version.
service discard
{
	disable		= no
	type		= INTERNAL
	id		= discard-stream
	socket_type	= stream
	protocol	= tcp
	user		= root
	wait		= no
}

# This is the udp version.
service discard
{
	disable		= no
	type		= INTERNAL
	id		= discard-dgram
	socket_type	= dgram
	protocol	= udp
	user		= root
	wait		= yes
}
